# docker/crontab
#execute sync.R every 5 minutes
5 * * * * Rscript /setup/sync.R
